# stage2/04-pirogue/files/pirogue_infos_screen.service
[Unit]
Description=PiRogue screen
Wants=influxdb.service

[Service]
ExecStart=/usr/bin/python3 /etc/pirogue/pirogue_infos_screen.py
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
Type=notify

[Install]
# Tell systemd to automatically start this service when the system boots
# (assuming the service is enabled)
WantedBy=default.target
