# stage2/04-pirogue/files/pirogue_daily_maintenance.timer
[Unit]
Description=Run PiRogue maintenance daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target