# stage2/04-pirogue/files/pirogue_daily_maintenance.service
[Unit]
Description=Run PiRogue maintenance

[Service]
Type=oneshot
ExecStart=/etc/pirogue/pirogue_maintenance.sh
ExecStop=
